# crontab
# émet l'événement EMI avec l'émission courante
* * * * * /var/www/overlay-engine/app/watch-show.js
