# crontab
# émet l'évenement EMI avec l'émission courante
*/5 * * * * /var/www/overlay-engine/app/push-show.js
